# templates/homeassistant.service.erb
# Managed with puppet.
[Unit]
Description=Home Assistant
After=network.target

[Service]
Type=simple
User=homeassistant
Environment=VIRTUAL_ENV="<%= @home %>"
Environment=PATH="<%= @home %>/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
ExecStart=<%= @home %>/bin/hass -c "<%= @confdir %>"

[Install]
WantedBy=multi-user.target

